# service/veyon.service.in
[Unit]
Description=Veyon Service
After=network-online.target dbus.service systemd-logind.service
Wants=network-online.target
Requires=dbus.service systemd-logind.service
Documentation=man:veyon-service(1)

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/veyon-service
Type=simple
Restart=always
StartLimitInterval=60
StartLimitBurst=10

[Install]
WantedBy=multi-user.target

